# systemd/cape.service
[Unit]
Description=CAPE
Documentation=https://github.com/kevoreilly/CAPEv2
Wants=cape-rooter.service
After=cape-rooter.service

[Service]
WorkingDirectory=/opt/CAPEv2/
ExecStart=/usr/bin/python3 cuckoo.py
#ExecStart=/opt/CAPEv2/venv/bin/python3 /opt/CAPEv2/cuckoo.py
User=cape
Group=cape
Restart=always
RestartSec=5m
LimitNOFILE=100000
[Install]
WantedBy=multi-user.target
