# systemd/suricata-update.service
[Unit]
Description=suricata-update

[Service]
User=root
Group=root
Type=oneshot
ExecStart=suricata-update

[Install]
WantedBy=multi-user.target
